# fb-crawler/scripts/launchd/ai.home4u.fb-chrome.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.home4u.fb-chrome</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_PATH__/fb-crawler/scripts/start-chrome.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <!-- Chrome immer am Leben halten. User-Cmd+Q triggert ebenfalls Restart;
       für Wartung den Job vorher entladen: launchctl bootout gui/$(id -u)/ai.home4u.fb-chrome -->
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/home4u/fb-chrome.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/home4u/fb-chrome.err.log</string>
</dict>
</plist>
